# src/buddy/micro-in-cn-tutorials/examples/basic-practices/micro-service/supervisord.conf
; ####################################
; @file: supervisord.conf
; @name: micro-service
; @date: 2019-12-18 12:35:58
; @tver: 20191217
; ####################################

[unix_http_server]
file=/tmp/go-stack-guide/sv.sock

[inet_http_server]
port=127.0.0.1:59002

[supervisord]
logfile=/tmp/go-stack-guide/sv.log
pidfile=/tmp/go-stack-guide/sv.pid
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
;nodaemon=true               ; start in foreground if true; default false


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix:///tmp/go-stack-guide/sv.sock
history_file=~/.sv_history


;[include]
;files = *.ini

[program:server-by-func]
command=make server-by-func-fg
stdout_logfile=/tmp/go-stack-guide/server-by-func-fg.log
redirect_stderr=true
stopasgroup=true
priority=8628

[program:server-by-service]
command=make server-by-service-fg
stdout_logfile=/tmp/go-stack-guide/server-by-service-fg.log
redirect_stderr=true
stopasgroup=true
priority=8628

[program:stream-file-server]
command=make start-stream-file-server-fg
stdout_logfile=/tmp/go-stack-guide/stream-file-server.log
redirect_stderr=true
stopasgroup=true
priority=8628

[program:stream-file-web]
command=make start-stream-file-web-fg
stdout_logfile=/tmp/go-stack-guide/stream-file-web.log
redirect_stderr=true
stopasgroup=true
priority=9999

[program:timeout-server]
command=make start-timeout-server-fg
stdout_logfile=/tmp/go-stack-guide/timeout-server.log
redirect_stderr=true
stopasgroup=true
priority=8628
